# deploy/ble-scale.service
[Unit]
Description=BLE Scale Reader
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=mikko
Group=bluetooth
WorkingDirectory=/home/mikko/ble-scale
EnvironmentFile=/home/mikko/ble-scale/.env
ExecStart=/home/mikko/ble-scale/ble-scale
Restart=always
RestartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ble-scale
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
